# home/private_dot_config/systemd/user/foot-server.socket
[Socket]
ListenStream=%t/foot.sock

[Unit]
ConditionEnvironment=WAYLAND_DISPLAY

[Install]
WantedBy=default.target
